# root/etc/services.d/transmission/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc transmission-daemon -f -g /bgmi/conf/transmission/ > /bgmi/log/transmission_out.log