# NAStool/etc/services.d/NAStool/run
#!/usr/bin/with-contenv bash

umask "${UMASK}"

exec s6-setuidgid nt /nas-tools