# root/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

exec \
  crond -f > /bgmi/log/crond.log