# dl_tools/aria2/bgmi_supervisord-aria2.ini
[unix_http_server]
file=%(ENV_BGMI_HOME)s/supervisor.sock
chmod=0700
[supervisord]
nodaemon=true
childlogdir=/bgmi/log/supervisor
logfile=%(ENV_BGMI_HOME)s/supervisord.log
pidfile=%(ENV_BGMI_HOME)s/supervisord.pid
user=root
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix://%(ENV_BGMI_HOME)s/supervisor.sock

[group:bgmi]
programs=bgmi_http,aria2,nginx,crond

[program:bgmi_http]
command=bgmi_http
user=bgmi
autorestart=true
stdout_logfile=/bgmi/log/bgmi_http.log
redirect_stderr=true

[program:aria2]
command=aria2c --conf-path=/bgmi/conf/aria2/aria2.conf
user=bgmi
autorestart=true
stdout_logfile=/bgmi/log/aria2_out.log
redirect_stderr=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
stdout_logfile=/bgmi/log/nginx_out.log
redirect_stderr=true

[program:crond]
command=crond -f
autorestart=true
stdout_logfile=/bgmi/log/crond.log
redirect_stderr=true