# dl_tools/aria2/bgmi_supervisord-aria2.ini
[group:bgmi]
programs=bgmi_http,aria2,nginx,crond

[program:bgmi_http]
command=bgmi_http
user=abc
autorestart=true
stdout_logfile=/bgmi/log/bgmi_http.log
redirect_stderr=true

[program:aria2]
command=aria2c --conf-path=/bgmi/conf/aria2/aria2.conf
user=abc
autorestart=true
stdout_logfile=/bgmi/log/aria2_out.log
redirect_stderr=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
stdout_logfile=/bgmi/log/nginx_out.log
redirect_stderr=true

[program:crond]
command=crond -f
autorestart=true
stdout_logfile=/bgmi/log/crond.log
redirect_stderr=true