# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec abc nginx -g "daemon off;"
