# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec \
  nginx -g "daemon off;" > /bgmi/log/nginx_out.log